# sdcard/rootfs/kubeadm/etc/systemd/system/docker.service.d/k8s.conf
[Service]
EnvironmentFile=/etc/kubernetes-on-arm/storagedriver.conf
ExecStart=
ExecStart=/usr/bin/dockerd -H fd:// --storage-driver=${STORAGE_DRIVER}
