# images/k8s/web/supervisord.conf
[program:http-server]
command=/bin/sh -c "cd /usr/lib/web/app; http-server -a 0.0.0.0 -p 8000"

[program:gulp]
command=/bin/sh -c "cd /usr/lib/web/master; npm start"